# arkdep-build.d/depends/arkanelinux-generic/overlay/post_bootstrap/etc/profile.d/abin.sh
# Utilized to prioritize wrapper applications in PATH
PATH="/usr/abin:$PATH"
